# config/bpfjailer-daemon.service
[Unit]
Description=BpfJailer Daemon (Standard Mode)
Documentation=https://github.com/gen0sec/jailer
After=network.target

[Service]
Type=simple
ExecStart=/usr/sbin/bpfjailer-daemon
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=no
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadOnlyPaths=/
ReadWritePaths=/sys/fs/bpf /run/bpfjailer

[Install]
WantedBy=multi-user.target
